# deploy/supercompute-backup.timer
[Unit]
Description=SuperCompute Daily Database Backup Timer

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target




